# services/healthcheck.service
[Unit]
Description=HTTP Health Check Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/healthcheck.sh
Restart=on-failure
RestartSec=5
SyslogIdentifier=healthcheck
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
